# deploy/systemd/diy-nav-backup.timer
[Unit]
Description=Back up diy-nav on the 1st and 16th of each month

[Timer]
OnCalendar=*-*-01,16 03:15:00
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
